# backend/scripts/com.tiergg.collect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tiergg.collect</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>/Users/johnnyseo/Desktop/master/Tier.GG/backend/scripts/collect_all.sh >> /tmp/tier.collect.log 2>&1</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/tier.collect.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/tier.collect.err</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>


